# Desktop/etc/systemd/system/pre-reboot-sign.path
[Unit]
Description=Watch for new OSTree deployments

[Path]
PathChanged=/ostree/deploy/$(hostname)/deploy
Unit=pre-reboot-sign.service

[Install]
WantedBy=multi-user.target
